# build/koel.conf
[program:koel]
user = nobody
command = /home/nobody/koel.sh
umask = 000

[program:mysql]
autorestart = false
startsecs = 0
user = root
command = /root/mysql.sh
umask = 000
