# setup/koel.conf
[program:koel]
user = nobody
command = /home/nobody/koel.sh

[program:mysql]
user = root
command = /root/mysql.sh
