# setup/koel.conf
[program:koel]
user = nobody
command = /home/nobody/koel.sh
umask = 000

[program:mysql]
user = root
command = /root/mysql.sh
umask = 000
